# bin/ubuntu/services/pianosite_search.service
[Unit]
Description=Pianosite Solr Search

[Service]
WorkingDirectory=/home/ubuntu/pianosite
Type=forking
KillMode=process
Restart=always
User=ubuntu
Group=ubuntu
ExecStart=/bin/bash -c 'cd /home/ubuntu/pianosite/solr-4.10.2/example/ && java -jar start.jar'

[Install]
WantedBy=multi-user.target
