# bin/ubuntu/services/pianosite_scheduler.service
[Unit]
Description=Pianosite Scheduler (Celery Beat)

[Service]
WorkingDirectory=/home/ubuntu/pianosite
User=ubuntu
Group=ubuntu
Environment=VIRTUAL_ENV=/home/ubuntu/.virtualenvs/pianosite
Environment=PATH=/home/ubuntu/.virtualenvs/pianosite/bin:$PATH
Environment=PATH=/usr/bin
ExecStart=/home/ubuntu/.virtualenvs/pianosite/bin/celery -A pianosite beat

[Install]
WantedBy=multi-user.target
