# bin/ubuntu/services/pianosite_scheduler.service
[Unit]
Description=Pianosite Scheduler

[Service]
WorkingDirectory=/home/ubuntu/pianosite
Type=forking
KillMode=process
Restart=always
User=ubuntu
Group=ubuntu
Environment=VIRTUAL_ENV=/home/ubuntu/.virtualenvs/pianosite
Environment=PATH=/home/ubuntu/.virtualenvs/pianosite/bin:$PATH
ExecStart=/home/ubuntu/.virtualenvs/pianosite/bin/celery -A pianosite beat

[Install]
WantedBy=multi-user.target
