# bin/ubuntu/services/pianosite_queue_monitor.service
[Unit]
Description=Pianosite Queue Monitor (Celery Flower)
Requires=pianosite_worker.service pianosite_scheduler.service

[Service]
WorkingDirectory=/home/ubuntu/pianosite
User=ubuntu
Group=ubuntu
Environment=VIRTUAL_ENV=/home/ubuntu/.virtualenvs/pianosite
Environment=PATH=/home/ubuntu/.virtualenvs/pianosite/bin:$PATH
Environment=PATH=/usr/bin
ExecStart=/home/ubuntu/.virtualenvs/pianosite/bin/flower -A pianosite --port=5555

[Install]
WantedBy=multi-user.target
